# config/launchagents/com.gchiam.tmux-theme-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gchiam.tmux-theme-watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gchiam/bin/tmux-theme-watcher</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/gchiam/.local/log/tmux-theme-watcher.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/gchiam/.local/log/tmux-theme-watcher.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/gchiam</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/gchiam</string>
    </dict>
</dict>
</plist>